# backup/supervisord.conf
[supervisord]

[program:backupdb]
command=/bin/backupdb